# .config/autostart/cloudapp.desktop
[Desktop Entry]
#Name=cloudapp
#Exec=/home/tyler/bin/cloudapp watch
#Type=Application
#Hidden=false

