# .config/autostart/droplr.desktop
[Desktop Entry]
Name=droplr
Exec=/home/tyler/bin/droplr
Type=Application
Hidden=false

