# .config/autostart/syncthing.desktop
[Desktop Entry]
Name=syncthing
Exec=syncthing
Type=Application
