# etc/systemd/undervolt.service
[Unit]
Description=undervolt

[Service]
Type=oneshot
ExecStart=/usr/bin/undervolt -v --core -120 --cache -120
